# minecraft-phantom/rootfs/etc/services.d/phantom/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting phantom..."

server=$(bashio::config 'server')
exec /usr/bin/phantom --server ${server}